# bin/halow/daemon/IoAire-halow.service
[Unit]
Description= Meta data provider of system details like version of v4l, mesh stats, wifi stats, and anything else needed
#Documentation="None available yet"

[Service]
ExecStart="/usr/bin/python3.11 /home/wilbur/.local/bin/halow/daemon/IoAire-daemon.py"
User=avahi

[Install]
WantedBy=multi-user.target

